# backup-etc/etc/cron.daily/fstrim
#!/bin/sh

test -x /sbin/fstrim || exit 0
/sbin/fstrim /
