# cron.daily
#!/bin/bash
# cron para rodar todo dia
/root/backup
