# backup-etc/etc/inittab
T0:23:respawn:/sbin/getty --autologin root -L ttyAMA0 115200 vt100
